# ansible/roles/chess-fracture/templates/etc/cron.d/purge_blend.j2

*/5 * * * * {{ ansible_user }} /bin/find /blend -type f -name "*.blend*" -mmin +30 -delete > /dev/null 2>&1
